# deploy/systemd/2busy1miss-runtime.timer
[Unit]
Description=2busy1miss Calendar Reminder Timer

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true

[Install]
WantedBy=timers.target
