# deploy/systemd/2busy1miss-runtime-agenda.timer
[Unit]
Description=2busy1miss Next-Day Calendar Agenda Timer

[Timer]
OnCalendar=*-*-* __AGENDA_SCHEDULE_TIME__:00
Persistent=true

[Install]
WantedBy=timers.target
