# deploy/systemd/2much2read-runtime.service
[Unit]
Description=2much2read Gmail Digest
After=network-online.target

[Service]
Type=oneshot
EnvironmentFile=%h/.config/2much2read-runtime/.2much2read.env
ExecStart=__EXECUTABLE__ run
TimeoutStartSec=30min
NoNewPrivileges=true
PrivateTmp=true
